# docker/rootfs/etc/services.d/app/run
#!/command/with-contenv sh
set -e

cd /app || exit 1
s6-setuidgid "$APP_USER" irc3 "$SETTINGS_FILE"
